# docker/grok/com.mc.grok-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mc.grok-bridge</string>
  <!-- Static example plist — replace /Users/YOUR_USER with your actual home dir before installing.
       ADR-063: grok is a headless host harness. The plist runs grok-bridge.py (the poll+dispatch
       HTTP server); each dispatch spawns a one-shot `grok` subprocess. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/YOUR_USER/Workspace/Projects/mission-control/scripts/grok-bridge.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <!-- KeepAlive: true — always restart, regardless of exit code (mirrors
       hermes-bridge Phase 26 HERM-12). grok-bridge.py raises SystemExit(1) on
       crash and exits 0 on SIGTERM, so launchd never restarts a clean stop. -->
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/YOUR_USER/.mc/agents/grok/logs/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USER/.mc/agents/grok/logs/launchd.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/YOUR_USER/.local/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/YOUR_USER</string>
  </dict>
</dict>
</plist>
